# packages/MainServiceAPI/mainserviceapi-0.0.4-py3-none-any.whl/MainServiceAPI/linux/blank.service
[Unit]
After=network.target
Description=Описание сервиса

[Service]
ExecStart=cat /dev/stdout
Restart=always
RestartSec=5
Type=simple
User=root
WorkingDirectory=/

[Install]
Alias=Название сервиса
WantedBy=multi-user.target